# FlyVPN_Gateway_Panel/systemd/blackwing-updater.timer
[Unit]
Description=Run BlackWing Gateway updater every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Unit=blackwing-updater.service
Persistent=true

[Install]
WantedBy=timers.target
