# FlyVPN_Gateway_Panel/systemd/flyvpn-panel.service
[Unit]
Description=FlyVPN Gateway Panel
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/flyvpn-gateway-panel
EnvironmentFile=/opt/flyvpn-gateway-panel/.env
ExecStart=/opt/flyvpn-gateway-panel/venv/bin/uvicorn app.main:app --host ${PANEL_HOST:-0.0.0.0} --port ${PANEL_PORT:-8090}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
